# chuwi-minibook-x-daemon/chuwi-minibook-x-daemon.service
[Unit]
Description=Chuwi Minibook X Daemon
Documentation=man:chuwi-minibook-x-daemon(8)
# Wait for IIO subsystem to be available
After=systemd-modules-load.service
# Wait for the kernel module to be loaded (if using modprobe)
After=sys-kernel-debug.mount
# Ensure the service starts after basic system initialization
After=sysinit.target

[Service]
Type=exec
ExecStart=/usr/sbin/chuwi-minibook-x-daemon
EnvironmentFile=-/etc/default/chuwi-minibook-x-daemon
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10

# Security settings
# Note: Must run as root to access kernel module sysfs and IIO devices
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Allow access to IIO sysfs and our kernel module sysfs
ReadOnlyPaths=/sys/bus/iio
ReadWritePaths=/sys/devices/platform/chuwi-minibook-x

[Install]
WantedBy=multi-user.target