# cmxsd/support/cmxsd-system.service
[Unit]
Description=CHUWI Minibook X Session Daemon (System-wide)
Documentation=man:cmxsd(8)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/cmxsd -c /etc/cmxsd/daemon.conf
Restart=on-failure
RestartSec=5
User=root

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Allow access to input devices and system resources
SupplementaryGroups=input

[Install]
WantedBy=multi-user.target