# accelerometers/99-mxc4005-second.rules
# Create the second MXC4005 (0x15) when the system is ready
# 
# This rule automatically creates the missing second accelerometer device
# on Chuwi Minibook X systems where ACPI/BIOS only exposes one device.

# Fire once when the ACPI parent appears (MDA6655:00 = \_SB_.ACMK)
ACTION=="add", SUBSYSTEM=="acpi", KERNEL=="MDA6655:00", RUN+="/usr/sbin/mda6655-split"

# Also fire when I2C adapters appear (covers resume/hotplug/races)
ACTION=="add", SUBSYSTEM=="i2c", KERNEL=="i2c-*", RUN+="/usr/sbin/mda6655-split"
